# supervisord.conf
[supervisord]

nodaemon=true

[program:gogs]

command=start.sh
